# openresty/openresty.conf
[program:openresty]
command=/opt/openresty/nginx/sbin/nginx -p /var/www -c /var/www/config/nginx.conf
autostart=true
